# deploy/image/stage-vanchor/01-net/files/vanchor-hotspot.service
[Unit]
Description=Vanchor setup-hotspot fallback
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/vanchor-hotspot-check

[Install]
WantedBy=multi-user.target
