# deploy/image/stage-vanchor/02-stack/files/var-log.mount
[Unit]
Description=Mount /var/log as tmpfs (SD-wear reduction)
DefaultDependencies=no
Before=local-fs.target

[Mount]
What=tmpfs
Where=/var/log
Type=tmpfs
Options=mode=0755,size=64M,noatime

[Install]
WantedBy=local-fs.target
